# Powershell/Task_scheduler_user_logoff_login/scheduler.xml
<!--This XML defines a Windows Scheduled Task that runs a specified PowerShell script automatically whenever any user logs off the system. 
The trigger is based on Security Event ID 4647, which corresponds to a user-initiated logoff event. 
The task runs with highest privileges under the SYSTEM account to ensure it has the necessary permissions to execute administrative actions, such as enabling network adapters or performing system-level logging. 
The script path specified in the <Arguments> section should be updated to point to the actual PowerShell script you wish to run at logoff. -->

<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Run a PowerShell script at user logoff (Event ID 4647)</Description>
    <Author>Administrator</Author>
  </RegistrationInfo>
  <Triggers>
    <EventTrigger>
      <Enabled>true</Enabled>
      <Subscription>
        <![CDATA[
          <QueryList>
            <Query Id="0" Path="Security">
              <Select Path="Security">*[System[EventID=4647]]</Select>
            </Query>
          </QueryList>
        ]]>
      </Subscription>
    </EventTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>SYSTEM</UserId>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-NoProfile -ExecutionPolicy Bypass -File "C:\Path\To\YourScript.ps1"</Arguments>
    </Exec>
  </Actions>
</Task>
